# run/etc/init.d/wifidog
#!/bin/sh /etc/rc.common

START=99

start() {
	wifidog
}

stop() {
	wdctl stop
}

reload() {
	wdctl restart
}
